# vagrant/virt/fake-ssds.service
[Unit]
Description=Fake SSDs

[Service]
Type=oneshot
ExecStart=/usr/bin/fake-ssds.sh

[Install]
WantedBy=multi-user.target
